# conf/supervisor/supervisord.ini
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/sbin/php-fpm7 --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'

[program:mikbotam-polling]
command=/usr/bin/php /public_html/Saldo/Core.php
directory=/public_html/Saldo
